# config/systemd/miniproject.service.d/goroot.conf
[Service]
Environment="GOPATH=/usr/local/lib/go"
